# packaging/deb/DEBIAN/postinst
#!/bin/sh

set -e

SERVICE_NAME="cc-switch-web.service"
STATE_DIR="/var/lib/cc-switch-web"
LOG_DIR="/var/log/cc-switch-web"

ensure_directories() {
    install -d -o root -g root "$STATE_DIR" "$LOG_DIR"
}

reload_and_start_service() {
    if ! command -v systemctl >/dev/null 2>&1; then
        return
    fi

    systemctl daemon-reload || true
    systemctl enable "$SERVICE_NAME" >/dev/null 2>&1 || true
    systemctl restart "$SERVICE_NAME" >/dev/null 2>&1 || systemctl start "$SERVICE_NAME" >/dev/null 2>&1 || true
}

case "$1" in
configure)
    ensure_directories
    reload_and_start_service
    ;;
esac

exit 0
